# examples/blockhost-monitor.service
[Unit]
Description=Blockhost Subscriptions Event Monitor
After=network.target

[Service]
Type=simple
EnvironmentFile=/opt/blockhost/.env
ExecStart=/usr/bin/blockhost-monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
